# config/service/papamonMasterApp.service
[Unit]
Description=PapamonMasterApp Python
[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/papamon/Documents/papamon/MasterApp.py
WorkingDirectory=/home/papamon/Documents/papamon
Restart=on-failure
RestartSec=10
[Install]
WantedBy=default.target
